# service/jarvis-camera.service
[Unit]
Description=Jarvis Camera Stream
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/rpicam-vid --inline --listen --signal=off -t 0 --width 640 --height 480 --framerate 15 --codec mjpeg -o tcp://0.0.0.0:8888
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
