# service/jarvis.service
[Unit]
Description=Jarvis Robot Control Server
After=network-online.target pigpiod.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/jarvis
ExecStart=/usr/local/bin/node /opt/jarvis/server.js
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
